# docker/cron-container/crontab
0 5 * * 1 curl http://web:8080/cron/sub
0 */12 * * * curl http://web:8080/cron/check_muted_vods
10 */12 * * * curl http://web:8080/cron/check_deleted_vods
